# files/etc/systemd/system/ftphealthcheck-restart@.service
[Unit]
Description=Restart service (%I)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart %i.service
